# packaging/hyprvoice.service
[Unit]
Description=Hyprvoice voice-to-text daemon
Documentation=https://github.com/leonardotrapani/hyprvoice
After=pipewire.service
Wants=pipewire.service

[Service]
Type=simple
ExecStart=/usr/bin/hyprvoice serve
Restart=on-failure
RestartSec=5
Environment=XDG_RUNTIME_DIR=/run/user/%i

[Install]
WantedBy=default.target